# systemd/polybar.service
Description=Autorun action for Polybar

[Service]
Type=oneshot
ExecStart=/usr/bin/su -c '/home/winson/.config/polybar/hdmi_plug.sh' winson
